# dot_config/systemd/user/updates-notifier.service
[Unit]
Description=Check and notify pacman updates

[Service]
Type=simple
ExecStart=/usr/bin/update-notifier
